# etc/crontabs/root
* * * * * lua /scripts/timeout.lua
* * * * * ../scripts/interstate.sh
*/5 * * * * lua /scripts/post-log.lua
*/5 * * * * lua /scripts/get-configuration.lua
