# access-point/etc/crontabs/root
* * * * * ../scripts/timeout.lua
*/5 * * * * lua /scripts/post-log.lua
